# motd.d/update-motd.d/update-motd.d
/home/adam/prometheus-dotfiles/motd.d/update-motd.d